# files/host/etc/systemd/system/getty@tty1.service.d/override.conf
# Override SUSE's override by clearing the TTY :)
# See /usr/lib/systemd/system/getty@tty1.service.d/20-defaults-SUSE.conf

[Service]
TTYVTDisallocate=yes
